# tennis_court_radar/rootfs/etc/s6-overlay/s6-rc.d/tennis-radar/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Tennis Court Radar..."
exec node /app/dist/bundle.js
