# investments/rootfs/etc/s6-overlay/s6-rc.d/investments/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Investments..."
exec node /app/dist/bundle.cjs
